# dotfiles/undervolt.service
[Unit]
Description=undervolt

[Service]
Type=oneshot
ExecStart=/usr/local/bin/undervolt -v --core -{{cpu_uv}} --gpu -{{gpu_uv}} --cache -{{cpu_uv}}
